# lib/systemd/podpourri-build-weekly@.timer
[Unit]
Description=Weekly container rebuilds

[Timer]
OnCalendar=weekly
RandomizedDelaySec=3600
Persistent=true

[Install]
WantedBy=timers.target
